# feature/builtin/roles/install/etcd/files/backup.service
[Unit]
Description=Backup ETCD
[Service]
Type=oneshot
ExecStart=/usr/local/bin/kube-scripts/backup_etcd.sh
